# packaging/jarvis.service
[Unit]
Description=JARVIS AI Voice Assistant
Documentation=https://github.com/YakupAtahanov/Project-JARVIS
After=network.target sound.target ollama.service
Wants=network.target ollama.service

[Service]
Type=simple
User=jarvis
Group=jarvis
# Hardware and system access groups needed for autonomous operation:
#   audio          — PipeWire microphone/speaker
#   video          — GPU/DRM access
#   network        — NetworkManager CLI socket
#   systemd-journal — journalctl log reading
#   storage        — block device access
SupplementaryGroups=audio video network systemd-journal storage

WorkingDirectory=/usr/lib/jarvis
RuntimeDirectory=jarvis
RuntimeDirectoryMode=0775
ExecStart=/usr/bin/jarvis-daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# CAP_NET_ADMIN allows direct network interface operations.
# CAP_SYS_NICE allows rtkit real-time scheduling for PipeWire audio.
# /dev/jarvis access is via group membership (GROUP=jarvis, MODE=0660) —
# CAP_SYS_ADMIN is not needed.
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_NICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_NICE

# NoNewPrivileges and RestrictSUIDSGID are intentionally absent:
# the daemon spawns `sudo` for system management (pacman, systemctl, etc.).
# sudo is a setuid binary — NO_NEW_PRIVS makes setuid a no-op, breaking it.
# The jarvis policy engine (SAFE/ELEVATED/DANGEROUS/FORBIDDEN) and the
# sudoers rules in /etc/sudoers.d/10-jarvis are the access-control layer.

# ProtectSystem is intentionally absent: package management (pacman) writes
# to /usr, /etc, and /var.  ProtectSystem=strict would create a read-only
# bind-mount namespace that blocks those writes even when running as root
# via sudo.

PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=JARVIS_CONFIG_DIR=/etc/jarvis
Environment=JARVIS_DATA_DIR=/var/lib/jarvis
Environment=JARVIS_INPUT_SOCKET=/run/jarvis/input.sock
Environment=JARVIS_OUTPUT_SOCKET=/run/jarvis/output.sock
Environment=JARVIS_LOG_DIR=/var/log/jarvis
Environment=JARVIS_MODELS_DIR=/var/lib/jarvis/models
Environment=PYTHONPATH=/usr/lib
Environment=OLLAMA_HOST=127.0.0.1:11434

# Standard output/error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jarvis

[Install]
WantedBy=multi-user.target
